# teamoon.service
[Unit]
Description=Teamoon - AI-powered project management and autopilot task engine
After=network.target

[Service]
Type=simple
User=cloud-agent
Group=cloud-agent
ExecStart=/usr/local/bin/teamoon serve
Restart=always
RestartSec=5
EnvironmentFile=/home/cloud-agent/Projects/teamoon/.env
WorkingDirectory=/home/cloud-agent

[Install]
WantedBy=multi-user.target
